# cron.daily/nncp-daily
#!/bin/bash

find /var/spool/nncp -name '*.seen' -a -mtime +30 -delete
find /var/spool/nncp -name '*.hdr' -a -mtime +7 -delete
